# www/dist/etc/podkstr.service
[Unit]
Description=podkstr
After=network-online.target
Wants=network-online.target systemd-networkd-wait-online.service

[Service]
Restart=on-failure
StartLimitInterval=86400
StartLimitBurst=5

; User and group the process will run as.
User=www-data
Group=www-data
WorkingDirectory=/var/www/podkstr.com

ExecStart=/var/www/podkstr.com/server
ExecReload=/bin/kill -USR1 $MAINPID

[Install]
WantedBy=multi-user.target